# etc/s6-overlay/s6-rc.d/minetestserver/finish
#!/bin/sh

service_exit_code=$1
uncaught_signal_code=$2

echo "$service_exit_code" > /run/s6-linux-init-container-results/exitcode
